# Script.cron
* * * * * echo "$(date) - Hourly log message" >> /home/eddy/makinates/cronjob-task/show.log
